# .config/systemd/user/rclone-sync.service
[Unit] 
Description=Sync keepass 

[Service] 
Type=simple 
ExecStart=/home/ernest/.scripts/rclone-sync.sh

[Install]
WantedBy=default.target
